# supervisor.conf
; ─── Regular short-lived jobs (campaigns, emails, drip, etc.) ────────────────
; retry_after on the 'database' connection is 90 s.  Keep --timeout below that
; (60 s) so a truly hung job is killed before the queue re-queues it.
[program:novelio-queue-worker]
command=php /var/www/html/bulk-email-platform/artisan queue:work database --queue=default --sleep=3 --tries=3 --timeout=60 --max-time=3600
directory=/var/www/html/bulk-email-platform
autostart=true
autorestart=true
startretries=5
numprocs=2
stopwaitsecs=70
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/novelio-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

; ─── Long-running CSV import jobs ────────────────────────────────────────────
; Uses the 'database_long' queue connection (retry_after: 1800 s).
; numprocs=1 ensures only ONE worker handles an import at a time — no second
; worker can ever reset the progress counter mid-run.
; --timeout=1800 must match (or exceed) the connection's retry_after.
[program:novelio-import-worker]
command=php /var/www/html/bulk-email-platform/artisan queue:work database_long --queue=long --sleep=3 --tries=2 --timeout=1800 --max-time=86400
directory=/var/www/html/bulk-email-platform
autostart=true
autorestart=true
startretries=5
numprocs=1
stopwaitsecs=1800
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/novelio-import-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
